# bash/ttyd.service
[Unit]
Description=ttyd
After=network.target


[Service]
ExecStart=/usr/local/ttyd -p 7681 zsh


[Install]
WantedBy=multi-user.target
